# OAP_sensor_data.service
[Unit]
Description=OAP Sensor Data Service
After=network-online.target

[Service]
Type=simple
User=pi
Group=pi
ExecStart=/usr/bin/python3 /usr/local/bin/oap-api/OAP_sensor_data.py
Restart=on-failure
RestartSec=5
StandardOutput=file:/home/pi/Logs/OAP_sensor_data.log
StandardError=file:/home/pi/Logs/OAP_sensor_data.error

[Install]
WantedBy=multi-user.target